# open-webui/rootfs/etc/s6-overlay/s6-rc.d/init-open-webui/up
#!/bin/bash
# ==============================================================================
# Home Assistant Add-on: Open WebUI
# Initializes Open WebUI configuration
# ==============================================================================

echo "[INFO] Initializing Open WebUI configuration..."

# Create data directory with proper permissions
if [[ ! -d "/data" ]]; then
    echo "[INFO] Creating data directory..."
    mkdir -p /data
    chown -R root:root /data
    chmod 755 /data
fi

# Create config directory if it doesn't exist
CONFIG_PATH="/data/config"
if [[ ! -d "${CONFIG_PATH}" ]]; then
    echo "[INFO] Creating config directory..."
    mkdir -p "${CONFIG_PATH}"
fi

# Generate secret key
echo "[INFO] Generating secret key..."
SECRET_KEY=$(python3 -c "import secrets; print(secrets.token_hex(32))")

# Set default values
OLLAMA_URL="${OLLAMA_BASE_URL:-http://supervisor:11434}"

# Create environment file
ENV_FILE="/data/.env"
echo "[INFO] Creating environment configuration..."
cat > "${ENV_FILE}" <<ENVEOF
# Open WebUI Configuration
DATA_DIR=/data
FRONTEND_BUILD_DIR=/app/build
STATIC_DIR=/app/backend/static

# Secret key for session management
WEBUI_SECRET_KEY=${SECRET_KEY}

# Ollama configuration
OLLAMA_BASE_URL=${OLLAMA_URL}

# Authentication settings
ENABLE_SIGNUP=true
ENABLE_LOGIN_FORM=true

# Web search configuration
ENABLE_WEB_SEARCH=false
WEB_SEARCH_ENGINE=searx

# Image generation configuration
ENABLE_IMAGE_GENERATION=false

# Network configuration
HOST=0.0.0.0
PORT=3000

# Database
DATABASE_URL=sqlite:///data/webui.db

# Model filtering
ENABLE_MODEL_FILTER=false
MODEL_FILTER_LIST=

# Default models
DEFAULT_MODELS=
ENVEOF

# Set permissions
chmod 600 "${ENV_FILE}"

# Initialize database directory
mkdir -p /data/db
chmod 755 /data/db

# Create logs directory
mkdir -p /data/logs
chmod 755 /data/logs

echo "[INFO] Open WebUI initialization completed!"
echo "[INFO] Secret key: ${SECRET_KEY:0:8}... (truncated for security)"
echo "[INFO] Ollama URL: ${OLLAMA_URL}"
